# templates/profile.d/slurm.csh.erb
#!/bin/csh
#
# File managed by Puppet!
#

setenv SLURM_CONF "<%= scope.lookupvar('slurm::slurm_conf_path') %>"
