# templates/profile.d/slurm.csh.erb
#!/bin/csh
#
# File managed by Puppet!
#

setenv SLURM_CONF="<%= @slurm_conf_path %>"
