# templates/profile.d/slurm.sh.erb
#!/bin/sh
#
# File managed by Puppet!
#

export SLURM_CONF="<%= scope.lookupvar('slurm::slurm_conf_path') %>"
