# templates/profile.d/slurm.csh.erb
#!/bin/csh

setenv SLURM_CONF="<%= @slurm_conf_path %>"
