# templates/profile.d/slurm.csh.erb
#!/bin/csh
#
# File managed by Puppet!
#

setenv SLURM_CONF "<%= scope['slurm::slurm_conf_path'] %>"
